# Raspberry_configuration/config_files/99-tty-raw.rules
KERNEL=="ttyAMA2", RUN+="/bin/stty -F /dev/ttyAMA2 raw -echo -echoe -echok -echoctl -echoke"
KERNEL=="ttyAMA3", RUN+="/bin/stty -F /dev/ttyAMA3 raw -echo -echoe -echok -echoctl -echoke"
